# infra/systemd/ziply5-web.service
[Unit]
Description=Ziply5 Next.js Web Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ziply5
Environment=NODE_ENV=production
Environment=PORT=5001
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
LimitNOFILE=65535
StandardOutput=append:/var/log/ziply5/web.log
StandardError=append:/var/log/ziply5/web-error.log

[Install]
WantedBy=multi-user.target
